# infrastructure/vms/bridges.service
[Unit]
Description=Cisco Live LTRSPG2212 – create host bridges for containerlab
Wants=network-pre.target
After=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ciscolive-ltrspg2212-create-bridges.sh
RemainAfterExit=yes

[Install]